# crontab
# Run batch pipeline daily at 6:00 AM UTC
0 6 * * * cd /app && node node_modules/tsx/dist/esm/cli.mjs scripts/batch.ts >> /proc/1/fd/1 2>&1
